# .config/systemd/user/vscode-prompt-sync.service
[Unit]
Description=Sync ~/notes/prompts/code.md -> VSCode custom.instructions

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync_code_to_vscode.sh
