# .config/systemd/user/wificheck.service
[Unit]
Description=Check if ISP fucked up

[Service]
Type=oneshot
ExecStart=%h/scripts/wificheck
