# .config/systemd/user/wificheck.timer
[Unit]
Description=Run wificheck periodically

[Timer]
OnBootSec=30s
OnUnitActiveSec=5s
Unit=wificheck.service

[Install]
WantedBy=timers.target
